# dotfiles/autostart/conky.desktop
[Desktop Entry]
Type=Application
Name=Conky
Exec=conky -b
X-GNOME-Autostart-enabled=true
